# taaos/init/services/taaos-neural-init.service
[Unit]
Description=TaaOS Neural Engine Initialization
Before=multi-user.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/taaos/init-neural-engine.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
